# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

chown -R root:root /etc/dropbear
chmod -R go-rwx /etc/dropbear

exit 0
